# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#


/usr/local/bin/wireless-connect.pl
cd /home/troy/rails && su troy -c "./start-both-servers.sh"
